# utilities/launchd/sites-watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mac9sb.sites-watcher</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/sh</string>
        <string>/Users/mac/Developer/utilities/scripts/sites-watcher.sh</string>
    </array>

    <key>WatchPaths</key>
    <array>
        <string>/Users/mac/Developer/sites</string>
    </array>

    <key>StartInterval</key>
    <integer>30</integer>

    <key>StandardOutPath</key>
    <string>/Users/mac/Library/Logs/com.mac9sb/sites-watcher.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/mac/Library/Logs/com.mac9sb/sites-watcher.stderr.log</string>
</dict>
</plist>
